# tokenpak/agent/systemd/tokenpak-proxy.service
[Unit]
Description=TokenPak LLM Proxy (Graceful Degradation)
Documentation=https://docs.tokenpak.dev/proxy
After=network.target
Wants=network.target

[Service]
Type=simple
# Replace %h with actual home path if running as a system service (not user service)
ExecStart=/usr/bin/env tokenpak serve
Restart=on-failure
RestartSec=5s
# Allow up to 5 restarts in 60 seconds before giving up
StartLimitIntervalSec=60
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tokenpak-proxy

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/.tokenpak

# Environment
Environment=PYTHONUNBUFFERED=1
# Add your API key env vars here (or use EnvironmentFile):
# EnvironmentFile=%h/.tokenpak/env

[Install]
# For user service: systemctl --user enable tokenpak-proxy
WantedBy=default.target
